# scheduler/systemd.service
# Linux systemd daily timer for newsagent.
# Install: copy both files to /etc/systemd/system/, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now newsagent.timer

# /etc/systemd/system/newsagent.service
[Unit]
Description=newsagent autonomous AI research report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/newsagent
ExecStart=/opt/newsagent/.venv/bin/python -m newsagent news /opt/newsagent/prompts/ai_news_daily.md
Environment=NEWSAGENT_CADENCE=daily
User=newsagent
